# files/systemd-binfmt.override.service
[Service]
ExecStartPre=/usr/bin/test -f /mnt/lima-rosetta/rosetta
RestartSec=5
Restart=on-failure
